# ops/systemd/1dentalai.service
[Unit]
Description=1DentalAI Next.js application
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/1DentalAI
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
